# debian/foldend.service
[Unit]
Description=Folden Service

[Service]
ExecStart=/usr/sbin/foldend run

[Install]
WantedBy=multi-user.target